# scripts/runit/mesh-status/run
#!/bin/sh
exec python3 -u /usr/local/bin/mesh-status.py