# scripts/runit/registrar/run
#!/bin/sh
exec python3 -u /usr/local/bin/registrar.py --loop
