# scripts/runit/sysinfo/run
#!/bin/sh
exec python3 -u /usr/local/bin/sysinfo.py --loop --interval 30