# scripts/runit/nginx/run
#!/bin/sh
printf '%s [nginx] service started\n' "$(date '+%Y-%m-%d %H:%M:%S %z')" >&2
exec nginx -g 'daemon off;'
